# dbserver/supervisord.conf
[supervisord]
nodaemon=true


[program:mysql]
command=/usr/sbin/mysqld
autorestart=true
autostart=true
